# config/puma.service
[Unit]
Description=Puma HTTP Server for carambus_api
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/carambus_api/current
Environment=PATH=/var/www/carambus_api/current/bin:/usr/local/bin:/usr/bin:/bin
Environment=RAILS_ENV=production
ExecStart=/var/www/carambus_api/current/bin/puma -C /var/www/carambus_api/shared/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always
RestartSec=1

# Ensure socket directory exists
ExecStartPre=/bin/mkdir -p /var/www/carambus_api/shared/sockets
ExecStartPre=/bin/mkdir -p /var/www/carambus_api/shared/pids
ExecStartPre=/bin/mkdir -p /var/www/carambus_api/shared/log

[Install]
WantedBy=multi-user.target
